# systemd/bluetty.service
[Unit]
Description=Bluetooth SPP to PTY Getty bridge
After=bluetooth.service dbus.service
Requires=bluetooth.service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStartPre=/usr/bin/test -f /etc/bluetty/bluetty.conf
ExecStart=/usr/local/bin/bluetty
Environment=RUST_LOG=info
Environment=BLUETTY_CONFIG=/etc/bluetty/bluetty.conf
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
